# ansible/roles/app_native/templates/barq.service.j2
[Unit]
Description=BARQ Lite Java Service
After=network-online.target
Wants=network-online.target

[Service]
User={{ app_user }}
Group={{ app_group }}
Environment=LOG_FILE={{ log_dir }}/barq.log
Environment=PORT={{ http_port }}
ExecStart=/usr/bin/java -jar {{ current_link }}/{{ jar_name }} --server.address=0.0.0.0 --server.port=8080

Restart=always
RestartSec=2

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ log_dir }}
MemoryMax=350M
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target

